# ansible/plays/roles/system/templates/fail2ban/update-blocklist.timer.j2
[Unit]
Description=Periodic update of fail2ban blocklist from Loki

[Timer]
# Run every hour
OnCalendar=hourly
# Add some randomness to prevent thundering herd
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
