# ansible/plays/roles/user/templates/oneshot.timer.j2
[Unit]
Description={{ timer.description }}

[Timer]
OnCalendar={{ timer.schedule }}
RandomizedDelaySec=5m
Persistent=true

[Install]
WantedBy=timers.target
