# ansible/plays/roles/ls_setup/templates/docker.mount.j2
[Unit]
Description=Bind mount Docker socket into jail
# Prevent systemd from auto-creating /run/docker.sock as a directory
# by ensuring docker.socket has already created the real socket file.
After=docker.socket
Requires=docker.socket
# Propagate docker.service stop/restart to this mount (one-way, no cycle).
# BindsTo would create a shutdown ordering cycle with local-fs.target.
PartOf=docker.service
# Skip default .mount ordering deps (local-fs.target, umount.target) that
# cause shutdown cycles with docker.service.
DefaultDependencies=no

[Mount]
What=/var/run/docker.sock
Where={{ ssh_chroot_jail_path }}/var/run/docker.sock
Type=none
Options=bind

[Install]
WantedBy=multi-user.target
